# files/streamdeck.path
[Unit]
Description="Stream Deck Device Path"

[Path]
# the device name will be different if you use streamdeck-mini or streamdeck-xl
PathExists=/dev/streamdeck-mini
Unit=streamdeck.service

[Install]
WantedBy=default.target